# agent/scripts/zedops-agent.service
[Unit]
Description=ZedOps Agent - Server Management Agent
Documentation=https://github.com/lasmarois/zedops
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/local/bin/zedops-agent --manager-url ${MANAGER_URL} --name ${AGENT_NAME}
Restart=always
RestartSec=10
TimeoutStopSec=30

# Run as root (required for Docker socket access)
User=root
Group=root

# Environment
Environment=HOME=/root
EnvironmentFile=-/etc/zedops-agent/config

# Security hardening (while maintaining Docker access)
NoNewPrivileges=no
ProtectSystem=false
ProtectHome=false
ReadWritePaths=/var/run/docker.sock /root/.zedops-agent

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zedops-agent

[Install]
WantedBy=multi-user.target
